# multiProcessContainer/supervisord.conf
[supervisord]
nodaemon=true
 
[program:tail]
command= tail -f /dev/null
autorestart=true
 
[program:sleep]
command= sleep infinity
autorestart=true